# build/filebotsetup.conf
[program:filebot-setup]
autorestart = false
startsecs = 0
user = root
command = /root/setupFileBot.sh
umask = 000